# www-apps/mastodon-openrc/files/mastodon-all
#!/sbin/openrc-run

. /etc/conf.d/mastodon

name="Mastodon meta service"

depend() {
    needs mastodon-web mastodon-streaming mastodon-sidekiq
}
